# deploy/agenthalo-dashboard.service
[Unit]
Description=AgentHALO Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agenthalo
Group=agenthalo
Environment=AGENTHALO_HOME=/var/lib/agenthalo
Environment=NUCLEUSDB_HOME=/var/lib/agenthalo
EnvironmentFile=-/etc/agenthalo/discord.env
ExecStart=/usr/local/bin/agenthalo dashboard --port 3100 --no-open
Restart=always
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agenthalo /etc/agenthalo

[Install]
WantedBy=multi-user.target
