# deploy/nucleusdb-p2p-bridge.service
[Unit]
Description=NucleusDB P2PCLAW Bridge Worker
After=network-online.target nucleusdb-dashboard.service
Wants=network-online.target nucleusdb-dashboard.service

[Service]
Type=simple
User=nucleusdb
Group=nucleusdb
WorkingDirectory=/opt/nucleusdb
Environment=AGENTHALO_HOME=/var/lib/nucleusdb
Environment=NUCLEUSDB_HOME=/var/lib/nucleusdb
EnvironmentFile=-/etc/nucleusdb/discord.env
EnvironmentFile=-/etc/nucleusdb/p2pclaw.env
ExecStart=/usr/bin/env agenthalo p2pclaw bridge run-loop --live --heartbeat --tau-sync
Restart=on-failure
RestartSec=15s
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/nucleusdb /etc/nucleusdb

[Install]
WantedBy=multi-user.target
